# services/incus-dns-sync/incus-dns-sync.service
[Unit]
Description=Incus DNS Sync — Lifecycle event dinleyici ve reconciler
Documentation=file:///opt/infra/services/incus-dns-sync/README.md
After=network-online.target incus.service
Wants=network-online.target
Requires=incus.service

[Service]
Type=simple
ExecStart=/usr/local/bin/incus-dns-sync -config /etc/incus-dns-sync/incus-dns-sync.yaml

# Beklenmedik kapanışlarda yeniden başlat
Restart=on-failure
RestartSec=10s

# Arka arkaya 5 başarısız denemeden sonra vazgeç
StartLimitIntervalSec=120s
StartLimitBurst=5

# Çıktı doğrudan journald'a
StandardOutput=journal
StandardError=journal
SyslogIdentifier=incus-dns-sync

# Güvenlik
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# Incus socket'e erişim gerekli
ReadWritePaths=/var/lib/incus

[Install]
WantedBy=multi-user.target
